# .config/autostart/gkrellm.desktop
[Desktop Entry]
Name=GKrellm
Comment=GKrellm Monitor 
Icon=gkrellm
Exec=gkrellm
Terminal=false
Type=Application
Hidden=true
