# .config/autostart/volumeicon.desktop
[Desktop Entry]
Name=Volume Icon
Comment=Volume Control for your system tray
Icon=
Exec=volumeicon
Terminal=false
Type=Application
Hidden=true
